# packaging/KoboRoot/etc/udev/rules.d/98-kscribbler.rules
KERNEL=="eth*", ACTION=="add", RUN+="/opt/bin/kscribbler"
KERNEL=="wlan*", ACTION=="add", RUN+="/opt/bin/kscribbler"
